# docker/crontab/root
* * * * * cd /var/www/html && /usr/local/bin/php artisan schedule:run >> /tmp/schedule.log 2>&1
